# rc.local
#!/bin/bash
if [[ -f "/etc/venus-os-expansions" ]]; then
    echo "venus-os-expansions already installed"
else
    bash /data/venus-os-expansions/install.sh
fi

